# grafana/supervisor/nginx.conf
[program:nginx]
command=nginx -g "daemon off;"
user=root
group=root
numprocs=1
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stdout
